# supervisor.conf
[supervisord]
nodaemon=true

[unix_http_server]
file=/var/run//supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run//supervisor.sock

[program:sshd]
command=/usr/sbin/sshd -D

[program:apache2]
priority=10
directory=/tmp
command=/usr/sbin/apache2ctl -D FOREGROUND
user=root
autostart=true
autorestart=true

[program:mariadb]
command=/usr/sbin/mysqld
priority=999
username=www-data
stdout_logfile=/tmp/%(program_name)s.stdout
stderr_logfile=/tmp/%(program_name)s.stderr

[program:postfix]
process_name=master
directory=/etc/postfix
command=/usr/sbin/postfix -c /etc/postfix start
startsecs= 0
autorestart=false

[program:tomcat]
environment=JAVA_HOME="/usr/lib/jvm/default-java",CATALINA_HOME="/usr/share/tomcat6",CATALINA_BASE="/var/lib/tomcat6",JAVA_OPTS="-Djava.awt.headless=true -Xmx128m -XX:+UseConcMarkSweepGC",CATALINA_PID="/var/run/tomcat6.pid",CATALINA_TMPDIR="/tmp/tomcat6-tomcat6-tmp",LANG="en_US.UTF-8",JSSE_HOME="/usr/lib/jvm/default-java/jre/"
command=/usr/share/tomcat6/bin/catalina.sh start
directory=/var/lib/tomcat6
priority=150
startsecs=10
user=tomcat6
group=tomcat6
startretries=0
stdout_logfile=/var/log/tomcat6/tomcat.log

[include]
files = /etc/supervisor/conf.d/*.conf